# Library/LaunchAgents/net.makerbox.dailycommand.plist
<?xml version="1.0" encoding="UTF-8"?><!-- -*- mode: xml; -*- -->
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>net.makerbox.dailycommand</string>
    <key>ProgramArguments</key>
    <array>
      <string>bash</string>
      <string>-c</string>
      <string>~/Library/LaunchAgents/daily.sh</string>
    </array>
    <key>StandardOutPath</key>
    <string>/tmp/daily.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/daily.err</string>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>43200</integer>
  </dict>
</plist>
